# ansible/roles/docker/files/profile.d/aliases.docker
docker() {
  local DOCKER_BINARY=$(which docker)
  case $1 in
    bash) shift; $DOCKER_BINARY exec -it "$1" /bin/bash ;;
    *) $DOCKER_BINARY "$@" ;;
  esac
}
