# ansible/roles/workstation/files/profile.d/aliases-git
gwip() { git add -A && git commit -m "${*:-Work in progress}" && git push; }
g() {
  case $1 in
    b) shift; git branch "$@" ;;
    c) shift; git checkout "$@" ;;
    l) git log --oneline -10 ;;
    s) git status ;;
    *) git "$@" ;;
  esac
}
