# ansible/roles/base/files/etc/cron.d/purge-logs
00 10 * * * find /var/log/* -type f -mtime +7 -exec rm {} \;